# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid junsu
chdir /home/junsu/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/superlists-staging.harveyk.me.socket \
    superlists.wsgi:application
